# auto-suspend/auto-suspend.service
[Unit]
Description=Auto-suspend VM after idle timeout
After=network.target

[Service]
Environment=IDLE_TIMEOUT=21600
ExecStart=/usr/local/bin/auto-suspend.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
