# auto-suspend/fix-network-on-resume.service
[Unit]
Description=Fix network interface after suspend/resume
After=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix-network-on-resume.sh

[Install]
WantedBy=suspend.target
